# systemd/dht-sensor.service
[Unit]
Description=MQTT DHT20 Sensor (Home Assistant)
After=network-online.target

[Service]
Type=simple
WorkingDirectory={{ PROJECT_DIR }}
ExecStart={{ PROJECT_DIR }}/.venv/bin/uv run python sensor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

