# systemd/lights.service
[Unit]
Description=MQTT Lights (Home Assistant)
After=network-online.target

[Service]
Type=simple
WorkingDirectory={{ PROJECT_DIR }}
ExecStart={{ PROJECT_DIR }}/.venv/bin/uv run python lights.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target

